# journal-scraper.service
[Unit]
Description=Journal Scraper Web Application
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/var/www/journal-scraper

# ---- Environment ----
Environment=FLASK_ENV=production
Environment=FLASK_DEBUG=0
Environment=SECRET_KEY=CHANGE_THIS_TO_A_LONG_RANDOM_SECRET
Environment=DATABASE_URL=postgresql://scraper_user:ScraperPass2024!@localhost:5432/journal_scraper
Environment=REDIS_URL=redis://localhost:6379/0
Environment=UPLOAD_FOLDER=/var/www/journal-scraper/results
Environment=DATA_FOLDER=/var/www/journal-scraper/data
Environment=PYTHONPATH=/var/www/journal-scraper

# ---- Process ----
ExecStart=/var/www/journal-scraper/venv/bin/gunicorn \
    --workers 2 \
    --threads 4 \
    --worker-class gthread \
    --timeout 120 \
    --bind 0.0.0.0:5000 \
    app:app

ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=journal-scraper

# ---- Security ----
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
